# scripts/launchd-clerk-morning.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.tradingagents.clerk.morning</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>exec /ABS/PATH/TradingAgents-main/scripts/cron-clerk-morning.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/ABS/PATH/TradingAgents-main</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>7</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/YOU/.tradingagents/logs/launchd-clerk-morning.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/YOU/.tradingagents/logs/launchd-clerk-morning.err.log</string>
</dict>
</plist>
